# systemd/waypointebooks.service
[Unit]
Description=waypoint_ebooks bot

[Service]
Type=oneshot
User=pi
WorkingDirectory=
ExecStart=python bot.py

[Install]
WantedBy=basic.target
